# tools/sca/examples/heartbeat/rootfs/etc/systemd/system/sca.service
[Unit]
Description=heartbeat (dstack self-contained app)
After=dstack-guest-agent.service

[Service]
ExecStart=/run/sca/bin/entrypoint.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
